# conf.d/novnc.conf
[program:novnc]
command=/home/novnc/repo-noVNC/utils/launch.sh --vnc localhost:5900 --listen 8080
autorestart=true
